# 00_システム/20_Agent_Portable/backup_automation/com.yuma.agent-runtime-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yuma.agent-runtime-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/kabushikikaishakorokoro/.codex/bin/agent_runtime_backup.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kabushikikaishakorokoro/.codex/bin/agent-runtime-live-backups/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kabushikikaishakorokoro/.codex/bin/agent-runtime-live-backups/launchd.err.log</string>
</dict>
</plist>
