# 01_プロジェクト/AI自動化/watchtower/com.korokoro.yuma-watchtower.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.korokoro.yuma-watchtower</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kojinn/.claude/scripts/run_watchtower.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/kojinn/.claude/scripts/watchtower.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kojinn/.claude/scripts/watchtower.err.log</string>
</dict>
</plist>
