# 00_システム/20_Agent_Portable/archive/claude/launchagents/com.claude.uber-weekly-plan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.uber-weekly-plan</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kojinn/.claude/scripts/run_uber_weekly_plan.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/kojinn/.claude/scripts/uber-weekly-plan-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kojinn/.claude/scripts/uber-weekly-plan-error.log</string>
</dict>
</plist>
