# 00_システム/20_Agent_Portable/archive/claude/launchagents/com.claude.discord-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.discord-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kojinn/.claude/scripts/discord-listener.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kojinn/.claude/scripts/discord-monitor-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kojinn/.claude/scripts/discord-monitor-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/kojinn</string>
        <key>PATH</key>
        <string>/Users/kojinn/.nvm/versions/node/v24.14.1/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>DISABLE_AUTOUPDATER</key>
        <string>1</string>
    </dict>
</dict>
</plist>
